# prepare_oim/roles/deploy_containers/pcs/templates/pcs-container.container.j2
[Unit]
Description=Pacemaker-Corosync HA Container for Podman-in-Podman

[Container]
ContainerName={{ pcs_container_name }}
Image={{ pcs_image_name }}:{{ pcs_image_tag }}
Volume=/etc/localtime:/etc/localtime:ro
Volume=/run/podman/podman.sock:/run/podman/podman.sock:ro
Volume={{ pcs_container_dir }}/pcs-start.sh:/opt/omnia/pcs/pcs-start.sh
Volume={{ pcs_container_dir }}/corosync:/etc/corosync{{ selinux_option }}
Volume={{ pcs_container_log_dir }}:/var/log/cluster{{ selinux_option }}
Volume={{ oim_shared_path }}/omnia:/opt/omnia{{ selinux_option }}
Network=host
PodmanArgs=--privileged

[Service]
Restart=always

[Install]
WantedBy=multi-user.target default.target